# zdenergy.service
[Unit]
Description=ZD Energy Sensor Dashboard
After=network.target

[Service]
Type=notify
# Replace 'your-username' with your actual Linux username on the droplet
User=your-username
WorkingDirectory=/home/your-username/zdenergy
Environment="PATH=/home/your-username/zdenergy/venv/bin"
; IMPORTANT: Use a single worker with threads so the TCP server (port 6000)
; runs exactly once per process and shares in-memory state with Flask routes.
; Multiple workers would cause port binding conflicts or split buffers.
ExecStart=/home/your-username/zdenergy/venv/bin/gunicorn -w 1 --threads 8 -b 0.0.0.0:5000 --timeout 120 readings:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
