# zdenergy.service
[Unit]
Description=ZD Energy Sensor Dashboard
After=network.target

[Service]
Type=notify
# Replace 'your-username' with your actual Linux username on the droplet
User=your-username
WorkingDirectory=/home/your-username/zdenergy
Environment="PATH=/home/your-username/zdenergy/venv/bin"
Environment="THREADS=8"
ExecStart=/home/your-username/zdenergy/venv/bin/gunicorn --workers 1 --threads ${THREADS} -b 0.0.0.0:5000 --timeout 120 readings:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
